# roles/server-common/templates/status-email-user@.service
[Unit]
Description=status email for %i to {{ ansible_user }}

[Service]
Type=oneshot
ExecStart=/usr/local/bin/systemd-email {{ ansible_user_email }} %i
User=nobody
Group=systemd-journal
