# roles/server-media/files/docker-compose.timer
[Unit]
Description=Pull and rebuild docker containers every week

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
