# roles/server-media/templates/aws-ddns@.service.j2
[Unit]
Description=Refresh Route53 DNS entry for %i
OnFailure=status-email-{{ ansible_user }}@%n.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/ddns-route53 -z {{ aws_ddns_hosted_zone }} -r %i
