# roles/server-media/templates/certbot.service.j2
[Unit]
Description=Let's Encrypt renewal
OnFailure=status-email-{{ ansible_user }}@%n.service

[Service]
Type=oneshot
ExecStart=/usr/bin/certbot renew -n --agree-tos
